# roles/update_scheduler/templates/wsl-setup.timer.j2
[Unit]
Description=Run wsl provisioning every 30 minutes

[Timer]
OnCalendar={{ provisioning_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
